# docker/root/etc/s6/giteria/run
#!/bin/bash
[[ -f ./setup ]] && source ./setup

pushd /app/giteria >/dev/null
exec su-exec $USER /usr/local/bin/giteria web
popd
